# config/includes.chroot/etc/systemd/system/mo-bodyd.service
[Unit]
Description=Morimil-controlled MO OS durable executor
Documentation=file:/usr/share/doc/mo-os/MORIMIL-EXECUTOR.md file:/usr/share/doc/mo-os/EXECUTOR-RECOVERY.md
After=local-fs.target
ConditionPathExists=/var/lib/mo-bodyd/identity/executor-private.pem
ConditionPathExists=/etc/mo/executor/pairing.json
ConditionPathExists=/etc/mo/executor/morimil-controller.pem

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/sbin/mo-executord serve
Restart=on-failure
RestartSec=2s
KillMode=control-group
TimeoutStopSec=30s
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/var/lib/mo-bodyd

[Install]
WantedBy=multi-user.target
