# config/includes.chroot/etc/systemd/system/mo-boot-ready.service
[Unit]
Description=Emit MO OS boot readiness signal
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/sbin/mo-boot-ready
RemainAfterExit=yes
