# config/includes.chroot/etc/systemd/system/mo-install-autotest.service
[Unit]
Description=Install MO OS onto an explicitly authorized QEMU disk
After=mo-boot-ready.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/sbin/mo-install-autotest run
TimeoutStartSec=20min
