# config/includes.chroot/etc/systemd/system/mo-installed-ready.service
[Unit]
Description=Emit MO OS installed-system readiness signal
After=multi-user.target
ConditionPathExists=/etc/mo-installed

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mo-installed-ready
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
