# config/includes.chroot/etc/systemd/system/mo-recovery-test-mutate.service
[Unit]
Description=Create one controlled post-snapshot mutation for recovery validation
After=local-fs.target
ConditionPathExists=/boot/mo-recovery-test-once

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/sbin/mo-recovery-test-mutate
TimeoutStartSec=2min

[Install]
WantedBy=multi-user.target
