# config/includes.chroot/etc/systemd/system/mo-recovery-verify.service
[Unit]
Description=Verify that the controlled mutation disappeared after rollback
After=local-fs.target
ConditionPathExists=/boot/mo-recovery-test-completed

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/sbin/mo-recovery-verify
TimeoutStartSec=2min

[Install]
WantedBy=multi-user.target
